# scripts/video-api.service
[Unit]
Description=Video Converter - API Server
Documentation=https://github.com/wingzc0/video-converter
After=network.target mariadb.service video-scanner.service video-processor.service
Wants=mariadb.service

[Service]
Type=simple
User={{SERVICE_USER}}
Group={{SERVICE_USER}}
WorkingDirectory={{INSTALL_DIR}}
EnvironmentFile={{INSTALL_DIR}}/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 {{INSTALL_DIR}}/daemon_ctl.py api start --foreground
Restart=on-failure
RestartSec=10
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=video-api

[Install]
WantedBy=multi-user.target
