# scripts/video-scanner.service
[Unit]
Description=Video Converter - Scan Daemon
Documentation=https://github.com/wingzc0/video-converter
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
User={{SERVICE_USER}}
Group={{SERVICE_USER}}
WorkingDirectory={{INSTALL_DIR}}
EnvironmentFile={{INSTALL_DIR}}/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 {{INSTALL_DIR}}/daemon_ctl.py scan start --foreground
ExecStop=/usr/bin/python3 {{INSTALL_DIR}}/daemon_ctl.py scan stop
Restart=on-failure
RestartSec=15
TimeoutStopSec=30
KillMode=mixed
StandardOutput=journal
StandardError=journal
SyslogIdentifier=video-scanner

[Install]
WantedBy=multi-user.target
